# root/crontab
0 2 * * * cd ~ && ./epg.sh > /dev/null 2>&1
